# scripts/crontab
# Refresh rotation every hour
0 */6 * * * ./scripts/refresh-rotation.sh
