# scripts/crontab
# Refresh rotation every hour
0 * * * * ./scripts/refresh-rotation.sh
